# etc/init.d/wpanet
#!/sbin/runscript
#service wpanet
source /etc/conf.d/wpanet

get_driver() (
    [ -e /sys/class/net/$1/wireless ] && printf "-Dwext " && exit
    [ -e /sys/class/net/$1/phy80211 ] && printf "-Dnl80211,wext " && exit
    case `readlink /sys/class/net/$1/device/driver` in
      */hostap*)
        printf  "-Dhostap "
	;;
      *)
        printf "-Dwired "
        ;;
    esac
)

start() {
    for IF in $IFACES
      do
            wpa_supplicant ${DRIVER:-`get_driver $IF`} -i $IF \
	      -c ${WPA_CFG:-/etc/wpa/$IF.conf} -B && \
            wpa_cli -i $IF -a /sbin/wpa_dhcp -B || \
	    echo "Could not start interface: $IF"
      done
}
stop() {
    for IF in $IFACES
      do
         /sbin/wpa_dhcp $IF stop
      done
}
help() {
    echo "Usage: $0 start | stop to start or stop all interfaces"
    echo "Use wpa_dhcp if0 stop to stop one interface"
}
